# commafeed/commafeed.service
[Unit]
Description=Commafeed (@DB)
After=network-online.target
Wants=network-online.target

[Service]
User=commafeed-@DB
Group=commafeed-@DB
WorkingDirectory=/var/lib/commafeed-@DB
ExecStart=/usr/bin/commafeed-@DB
SuccessExitStatus=143

[Install]
WantedBy=multi-user.target
